# infra/pi/kpp-port5000-proxy.service
[Unit]
Description=Proxy legacy KPP port 5000 to FastAPI port 8000
After=kpp-dashboard.service

[Service]
ExecStart=/usr/bin/python3 /home/evan/kpp/infra/pi/port5000_proxy.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
